# tiddly/rootfs/etc/services.d/tiddly/run
#!/command/with-contenv bashio

PORT=$(bashio::config 'port')
DB_PATH=$(bashio::config 'db_path')
PLAINTEXT_PASSWORD=$(bashio::config 'admin_password')

mkdir -p "$(dirname "$DB_PATH")"

# Optional: Passwort-Hash erzeugen, falls gesetzt
if bashio::var.has_value "$PLAINTEXT_PASSWORD"; then
    echo "[INFO] Plaintext password provided, generating bcrypt hash..."
    HASH=$(deno --allow-env eval 'import { hashSync } from "https://deno.land/x/bcrypt/mod.ts"; console.log(hashSync(Deno.env.get("PLAINTEXT_PASSWORD") ?? "admin", 10));')
    export ADMIN_PASSWORD_HASH="$HASH"
    echo "[INFO] Generated admin password hash is ${#HASH} characters long"
else
    echo "[INFO] No admin password set. Open access."
    export ADMIN_PASSWORD_HASH=""
fi

export DB_PATH="$DB_PATH"

# Start server
exec deno run --unstable-broadcast-channel \
  --allow-env \
  --allow-read="$(dirname "$DB_PATH")" \
  --allow-write="$(dirname "$DB_PATH")" \
  --allow-net=":${PORT}" \
  https://codeberg.org/valpackett/tiddlypwa/raw/branch/release/server/run.ts
