# tiddly/rootfs/services.d/tiddlywiki/run
#!/command/with-contenv bashio

PORT=$(bashio::config 'port')
DB_PATH=$(bashio::config 'db_path')
HASH=$(bashio::config 'admin_password_hash')
SALT=$(bashio::config 'admin_password_salt')

mkdir -p "$(dirname "$DB_PATH")"

export ADMIN_PASSWORD_HASH=$HASH
export ADMIN_PASSWORD_SALT=$SALT
export DB_PATH=$DB_PATH

exec deno run --unstable-broadcast-channel \
  --allow-env \
  --allow-read="$(dirname "$DB_PATH")" \
  --allow-write="$(dirname "$DB_PATH")" \
  --allow-net=":${PORT}" \
  https://codeberg.org/valpackett/tiddlypwa/raw/branch/release/server/run.ts
