# tiddly/rootfs/etc/services.d/caddy/run
#!/usr/bin/with-contenv bash
exec caddy run --config /etc/caddy/Caddyfile --adapter caddyfile