# tiddly/rootfs/etc/services.d/caddy/run
#!/command/with-contenv bashio

export CERTFILE=$(bashio::config 'certfile')
export KEYFILE=$(bashio::config 'keyfile')

envsubst < /etc/caddy/Caddyfile.template > /etc/caddy/Caddyfile

ls -l /ssl

cat /etc/caddy/Caddyfile

if [ ! -s /etc/caddy/Caddyfile ]; then
  echo "❌ Error: /etc/caddy/Caddyfile is empty or invalid!"
  exit 1
fi

exec caddy run --config /etc/caddy/Caddyfile --adapter caddyfile