# tiddly/rootfs/etc/services.d/caddy/run
#!/command/with-contenv bashio

export CERTFILE=$(bashio::config 'certfile')
export KEYFILE=$(bashio::config 'keyfile')

envsubst < /etc/caddy/Caddyfile.template > /etc/caddy/Caddyfile

ls -l /ssl

cat /etc/caddy/Caddyfile

if [ ! -s /etc/caddy/Caddyfile ]; then
  echo "❌ Error: /etc/caddy/Caddyfile is empty or invalid!"
  exit 1
fi

exec caddy run --config /etc/caddy/Caddyfile --adapter caddyfile

# Start Caddy
caddy start --config /etc/caddy/Caddyfile

# Warte kurz, bis der Dienst hochgefahren ist
sleep 2

# Prüfen, ob Port 8443 geöffnet ist
if nc -z localhost 8443; then
  echo "✅ Caddy mit SSL aktiv auf Port 8443"
else
  echo "❌ Caddy/SSL nicht aktiv auf Port 8443"
fi
