# tiddly/rootfs/etc/services.d/caddy/run
#!/command/with-contenv bashio

export CERTFILE=$(bashio::config 'certfile')
export KEYFILE=$(bashio::config 'keyfile')

envsubst < /etc/caddy/Caddyfile.template > /etc/caddy/Caddyfile

ls -l /ssl

cat /etc/caddy/Caddyfile

if [ ! -s /etc/caddy/Caddyfile ]; then
  echo "❌ Error: /etc/caddy/Caddyfile is empty or invalid!"
  exit 1
fi

caddy start --config /etc/caddy/Caddyfile --adapter caddyfile
sleep 2

if nc -z localhost 8443; then
  echo "✅ Caddy mit SSL aktiv auf Port 8443"
else
  echo "❌ Caddy/SSL nicht aktiv auf Port 8443"
  exit 1
fi

# Stoppe Caddy Background-Instance, wenn du stattdessen `run` verwenden willst
caddy stop

# Starte Caddy im Vordergrund für HA
exec caddy run --config /etc/caddy/Caddyfile --adapter caddyfile